# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=warn
pidfile=/run/supervisord.pid

# PHP-FPM starts first (priority 1)
[program:php-fpm]
command=php-fpm -F
priority=1
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Nginx starts second (priority 2)
[program:nginx]
command=nginx -g 'daemon off;'
priority=2
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Init runs AFTER fpm+nginx are up (priority 10).
# Handles: package discover, config cache, DB wait, migrations, first-install seed.
# autorestart=false — only runs once per container lifecycle.
# TEMPORARILY DISABLED for debugging - uncomment when healthcheck works
#[program:init]
#command=/bin/sh -c "sleep 5 && /init.sh"
#priority=10
#autostart=true
#autorestart=false
#startretries=0
#startsecs=0
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
